# crontab
## Вычисление хэшей файлов
0 * * * *       php /path/to/generate_all_hashes.php > /dev/null
## Обновление iso
0 0 * * 6       php /path/to/update_iso.php > /dev/null
